# bin/luxcena-neo.service
[Unit]
Description=Luxcena Neo

[Service]
ExecStart=/opt/luxcena-neo/bin/luxcena-neo.sh
Restart=always
RestartSec=10

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=development
WorkingDirectory=/opt/luxcena-neo/

[Install]
WantedBy=multi-user.target
