# bin/luxcena-neo.service
[Unit]
Description=Luxcena Neo

[Service]
ExecStart=/home/lux-neo/src/bin/luxcena-neo.sh
Restart=always
RestartSec=10

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
WorkingDirectory=/home/lux-neo/install/src/

[Install]
WantedBy=multi-user.target
