# bin/luxcena-neo.service
[Unit]
Description=Luxcena Neo

[Service]
ExecStart={{WD}}/bin/luxcena-neo.sh
Restart=always
RestartSec=10

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=development
WorkingDirectory={{WD}}

[Install]
WantedBy=multi-user.target
